# 2-processes-and-services/scripts/sample-worker.service
[Unit]
Description=Sample Worker Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /scripts/sample-worker.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
